# package/base-files/files/etc/uci-defaults/16_prplos_cron_logrotate
grep -s -q /usr/sbin/logrotate /etc/crontabs/root && exit 0
echo "# logrotate runs every 10 minutes" >> /etc/crontabs/root
echo "*/10 * * * * /usr/sbin/logrotate /etc/logrotate.conf" >> /etc/crontabs/root
exit 0
