# package/base-files/files/etc/uci-defaults/20_prplos_no_coredumps
is_likely_tr181_component_with_coredumps_enabled()
{
	local init_script="$1"
	grep -sq -e 'ulimit -c unlimited' -e 'ubus-cli' "$init_script"
}

for init_script in /etc/init.d/*; do
	if is_likely_tr181_component_with_coredumps_enabled "$init_script"; then
		logger -t "uci-defaults" -p daemon.warn "prplOS: ${init_script}: disabling 'ulimit -c unlimited'"
		sed -i 's/ulimit -c unlimited/# ulimit -c unlimited/g' "$init_script"
	fi
done

exit 0
